# roles/node-exporter/files/docker.timer
[Unit]
Description=Runs the docker.service 1 minute after boot up

[Timer]
OnBootSec=1 m
Unit=docker.service

[Install]
WantedBy=basic.target